# deploy/meta-search.service
[Unit]
Description=Meta Search MCP Server
Documentation=https://github.com/lieyan666/mcp-server-meta-search
After=network.target

[Service]
Type=simple
User=meta-search
Group=meta-search
WorkingDirectory=/opt/meta-search
ExecStart=/opt/meta-search/meta-search
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=CONFIG_PATH=/opt/meta-search/config.jsonc
Environment=LOGS_DB_PATH=/opt/meta-search/data/logs.db
Environment=PORT=3000
Environment=HOST=127.0.0.1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/meta-search/config.jsonc /opt/meta-search/data
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meta-search

[Install]
WantedBy=multi-user.target
